# dot_config/systemd/user/xppen.service.tmpl
[Unit]
Description=XP-Pen Initialization Script

[Service]
Type=oneshot
ExecStart={{ .scriptsDir }}/xppen-attach
RemainAfterExit=yes

[Install]
WantedBy=default.target   

# vim: set ft=systemd :
