# dot_config/profile.d/flux
# completion
if command -v flux &> /dev/null; then
    . <(flux completion zsh)
fi
