# dot_config/profile.d/gitpod
# completion
# FIXME
# craaaazy slow!
# if command -v gitpod &> /dev/null; then
#     . <(gitpod completion zsh)
# fi
