# dot_config/profile.d/gitpod
# completion
if command -v gitpod &> /dev/null; then
    . <(gitpod completion zsh)
fi
