# dot_config/profile.d/navi
if command -v navi &> /dev/null; then
    _evalcache navi widget zsh
fi
