# dot_config/profile.d/pandoc
export TEXINPUTS=${HOME}/.config/pandoc/themes/dracula//:
alias pandf="pandoc --lua-filter ${HOME}/.pandoc/ext/diagram/diagram.lua --pdf-engine=xelatex --defaults ${HOME}/.config/pandoc/themes/dracula/dracula.yaml -V mainfont=\"FiraCode Nerd Font Mono\""
