# dot_config/profile.d/pandoc
alias pandf="pandoc --lua-filter ${HOME}/.pandoc/ext/diagram/diagram.lua --pdf-engine=xelatex"
