# dot_config/profile.d/path
# Function to add directory to PATH if it exists
add_to_path() {
    local dir="$1"
    if [ -d "$dir" ]; then
        export PATH="${PATH}:${dir}"
    fi
}

# Function to prepend directory to PATH if it exists (for priority)
prepend_to_path() {
    local dir="$1"
    if [ -d "$dir" ]; then
        export PATH="${dir}:${PATH}"
    fi
}

# Add directories to PATH (appended)
add_to_path "${HOME}/.local/bin"

# this needs to be set before .zshrc's plugins=(...)
#add_to_path "${HOME}/.local/scripts"

add_to_path "${HOME}/bin"
add_to_path "${HOME}/scripts"
add_to_path "${HOME}/android-studio/bin"

# Prepend directories to PATH (for priority)
prepend_to_path "/opt/homebrew/opt/mysql-client/bin"
add_to_path "${HOME}/.cache/.bun/bin"
