# dot_config/profile.d/pyenv
# export PYENV_ROOT="$HOME/.pyenv"
# [[ -d $PYENV_ROOT/bin ]] && export PATH="${PYENV_ROOT}/bin:$PATH"
# if command -v pyenv &> /dev/null; then
#     # FIXME - super slow
#     # _evalcache pyenv init -
#     # _evalcache pyenv virtualenv-init -
# fi
