# dot_config/profile.d/xargs
# get xargs to use aliases

alias xargs='xargs '
