# dot_config/profile.d/zephyr
if [ -f "${HOME}/.zephyrrc" ]; then source "${HOME}/.zephyrrc"; fi

