# dot_config/profile.d/zoxide
if command -v zoxide &> /dev/null; then
    _evalcache zoxide init zsh
    alias cd="z"
fi
